# cdr_mtn_tv/scripts/systemd/cdr-mtn-tv-displays.service
[Unit]
Description=cdr_mtn_tv feh display viewers (TV1 + TV2)
After=graphical.target cdr-mtn-tv-web.service
Wants=cdr-mtn-tv-web.service

[Service]
Type=forking
User=__USER__
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/scripts/start_displays.sh
ExecStop=/usr/bin/pkill -f "feh.*tv1_menu.jpg"; /usr/bin/pkill -f "feh.*tv2_events.jpg"
Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical.target
