# cdr_mtn_tv/scripts/systemd/cdr-mtn-tv-refresh.service
[Unit]
Description=cdr_mtn_tv startup events refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/venv/bin/python scripts/refresh_events.py

[Install]
WantedBy=multi-user.target
